# deploy/retired/polybot-strategy-g.service
[Unit]
Description=polyBOT strategy paper bot — variant G (BTC 5m, obs-3, thr-3, fill<=0.92, blackout 13/17/22). Phase 3 redirect to 5m markets.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polybot
Group=polybot
WorkingDirectory=/opt/polybot/polymarket-cli
Environment=STRATEGY_DATA_DIR=/opt/polybot/polymarket-cli/scripts/strategy/data-g
Environment=STRATEGY_SLOT_SECS=300
Environment=STRATEGY_SLUG_PREFIX=btc-updown-5m-
Environment=MAX_FILL_PRICE=0.92
Environment=STRATEGY_BLACKOUT_HOURS=13,17,22
ExecStart=/usr/bin/node /opt/polybot/polymarket-cli/scripts/strategy/main.js 3 3 100 168
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot-strategy-g
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/polybot/polymarket-cli/scripts/strategy/data-g

[Install]
WantedBy=multi-user.target
