# deploy/polybot-strategy-d.service
[Unit]
Description=polyBOT strategy paper bot — variant D (B refined: observe-11, threshold-6bps, max-fill-0.92, $100/trade). Replaces C (10bps) which produced 0 entries in 11h.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polybot
Group=polybot
WorkingDirectory=/opt/polybot/polymarket-cli
Environment=STRATEGY_DATA_DIR=/opt/polybot/polymarket-cli/scripts/strategy/data-d
Environment=MAX_FILL_PRICE=0.92
ExecStart=/usr/bin/node /opt/polybot/polymarket-cli/scripts/strategy/main.js 11 6 100 168
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot-strategy-d
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/polybot/polymarket-cli/scripts/strategy/data-d

[Install]
WantedBy=multi-user.target
