# deploy/polybot-strategy-k.service
[Unit]
Description=polyBOT strategy paper bot — variant K (obs-11, thr-6, fill<=0.92, blackout 13/17/22 UTC, stopLoss=10bps reversal). BTC 15m + stop-loss test vs D.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polybot
Group=polybot
WorkingDirectory=/opt/polybot/polymarket-cli
Environment=STRATEGY_DATA_DIR=/opt/polybot/polymarket-cli/scripts/strategy/data-k
Environment=MAX_FILL_PRICE=0.92
Environment=STRATEGY_BLACKOUT_HOURS=13,17,22
Environment=STOP_LOSS_RETBPS_REVERSAL=10
ExecStart=/usr/bin/node /opt/polybot/polymarket-cli/scripts/strategy/main.js 11 6 100 168
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot-strategy-k
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/polybot/polymarket-cli/scripts/strategy/data-k

[Install]
WantedBy=multi-user.target
