# deploy/polybot-strategy-h.service
[Unit]
Description=polyBOT strategy paper bot — variant H (ETH 15m, obs-11, thr-8, fill<=0.92, blackout 13/17/22). Phase 3 redirect to ETH.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polybot
Group=polybot
WorkingDirectory=/opt/polybot/polymarket-cli
Environment=STRATEGY_DATA_DIR=/opt/polybot/polymarket-cli/scripts/strategy/data-h
Environment=STRATEGY_SLUG_PREFIX=eth-updown-15m-
Environment=STRATEGY_PRICE_MODULE=./eth_price
Environment=MAX_FILL_PRICE=0.92
Environment=STRATEGY_BLACKOUT_HOURS=13,17,22
ExecStart=/usr/bin/node /opt/polybot/polymarket-cli/scripts/strategy/main.js 11 8 100 168
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=polybot-strategy-h
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/polybot/polymarket-cli/scripts/strategy/data-h

[Install]
WantedBy=multi-user.target
