# scripts/ui/__fixtures__/deploy/polybot-strategy-malformed.service
[Unit]
Description=intentionally broken - missing ExecStart
[Service]
User=polybot
[Install]
WantedBy=multi-user.target
